# debian/monitoring/systemd/socket-snoop.service
[Unit]
Description=Socket Snoop eBPF monitor
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
# Set this to your repo root (the Makefile's service-install auto-fills it):
# Example after sed: WorkingDirectory=/root/linux-tools/debian/monitoring
WorkingDirectory=/opt/socket-snoop/debian/monitoring
ExecStart=/usr/bin/python3 /opt/socket-snoop/debian/monitoring/socket_snoop.py --log-file=/var/log/socket_monitor.log
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
