# debian/monitoring/systemd/socket-snoop.service
# /etc/systemd/system/socket-snoop.service
[Unit]
Description=Socket Snoop eBPF monitor
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/linux-tools/debian/monitoring
ExecStart=/usr/bin/python3 /root/linux-tools/debian/monitoring/socket_snoop.py --log-file=/var/log/socket_monitor.log
Restart=on-failure
RestartSec=3

# Optional hardening (root still required for eBPF; drop if it causes issues)
# CapabilityBoundingSet=CAP_BPF CAP_SYS_ADMIN CAP_NET_ADMIN
# AmbientCapabilities=CAP_BPF CAP_SYS_ADMIN CAP_NET_ADMIN
# NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
