# debian/monitoring/systemd/socket-snoop.service
[Unit]
Description=Socket Snoop eBPF monitor
After=network.target
Wants=network.target

[Service]
Type=simple
# Replace /opt/socket-snoop with your absolute repo path
ExecStart=/usr/bin/env bash -lc '/opt/socket-snoop/.venv/bin/python /opt/socket-snoop/debian/monitoring/socket_snoop.py --log-file=/var/log/socket_monitor.log'
User=root
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
